# supervisor/GPIO-Slave.conf
[program:GPIO-Slave]
command = sudo nohup python ./GPIO-Slave.py &
directory = /home/pi/gpio/GPIO-Web
user = pi