# supervisor/GPIO-Master.conf
[program:GPIO-Master]
command = sudo nohup python ./GPIO-Master.py &
directory = /home/pi/gpio/GPIO-Web
user = pi